# vless-generator.service
[Unit]
Description=VLESS Config Generator Service
Documentation=https://github.com/aaleksandrov/vless-generator
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/vless-generator
ExecStart=/opt/vless-generator/vless-generator \
    -port=8080 \
    -log-level=info \
    -log-format=json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-generator

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
