# vless-generator.service
[Unit]
Description=VLESS Config Generator Service
Documentation=https://github.com/aaleksandrov/vless-generator
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/vless-generator
ExecStart=/opt/vless-generator/vless-generator \
    -port=8080 \
    -server=vless.example.com \
    -server-port=443 \
    -ws-path=/websocket \
    -dns-server=8.8.8.8 \
    -doh-server=https://223.5.5.5/dns-query \
    -tun-address=172.19.0.1/28 \
    -mixed-port=2080 \
    -tun-mtu=9000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-generator

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
